# deploy/nodeble.service
[Unit]
Description=NODEBLE Iron Condor Bot
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/nodeble
ExecStart=%h/nodeble/.venv/bin/python -m nodeble --mode scan
Environment=PYTHONPATH=%h/nodeble/src

[Install]
WantedBy=multi-user.target
